# kernel-extension/templates/systemd/oxisentinel.service.d/10-limits.conf
[Service]
# Placeholder for OxiSentinel service limits.
